# notify.service
[Unit]
Description=Run notification task every hour

[Service]
Type=simple
ExecStart=/bin/notify

[Install]
WantedBy=multi-user.target